# caddy.ini
[program:caddy]
command=/bin/caddy -conf /etc/caddy/caddy.conf -http2 -quic
stderr_logfile = /var/log/caddy-err.log
stdout_logfile = /var/log/caddy-log.log