# nginx.ini
[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
stderr_logfile = /var/log/nginx-err.log
stdout_logfile = /var/log/nginx-log.log