# php.ini
[program:php]
command=php-fpm --nodaemonize
stderr_logfile = /var/log/php-err.log
stdout_logfile = /var/log/php-log.log