# varnish.ini
[program:varnish]
command=/usr/bin/varnishd -a :77 -F -f /etc/varnish/default.vcl -T 127.0.0.1:6082 -p thread_pool_min=50 -p thread_pool_max=1000 -S /etc/varnish/secret -s malloc,256M -P /var/run/varnish.pid
stderr_logfile = /var/log/varnish-err.log
stdout_logfile = /var/log/varnish-log.log